# systemd/proxy-tunnel.target
[Unit]
Description=Tunnel outgoing TCP/IPv4 traffic through a proxy
Documentation=https://github.com/blokyk/ebpf-tunnel